# sv/syncthing/run
#!/bin/execlineb
fdmove -c 1 2
syncthing serve --no-browser
