# sv/xdg-desktop-portal-gtk/run
#!/bin/execlineb
/usr/libexec/xdg-desktop-portal-gtk
