# sv/spawnd/run
#!/bin/sh
mkfifo "$SPAWND_FIFO" || true
printf %s "$HOME" >"$SPAWND_FIFO"
exec spawnd
