# sv/spawnd/run
#!/bin/sh
cd ~
echo ~ > "$XDG_RUNTIME_DIR/spawnpwd"
exec spawnd
