# sv/xdg-desktop-portal/run
#!/bin/execlineb
/usr/libexec/xdg-desktop-portal -r
