# sv/xdg-desktop-portal-wlr/run
#!/bin/execlineb
/usr/libexec/xdg-desktop-portal-wlr
