# sv/execd/run
#!/bin/sh
mkfifo "$EXECD_FIFO" || true
changepwd "$HOME" &
exec execd
