# sv/rivertile/run
#!/bin/sh
exec rivertile -view-padding 10 -outer-padding 5
