# sv/wbg/run
#!/bin/execlineb
importas -S HOME
wbg -s ${HOME}/Pictures/wallpapers/wallpaper.png
