# sv/wbg/run
#!/usr/bin/execlineb
importas -S HOME
wbg ${HOME}/Pictures/wallpapers/ashen/dark-souls.png
