# sv/transmission/run
#!/bin/execlineb
transmission-daemon -f
